# modules/supernova/files/create.cron
MAILTO=sm@ocf.berkeley.edu
# m h dom mon dow user   command
  42 16 *   *   *   create /opt/ocf/packages/create/create.py --keytab=/opt/create/private/create.keytab

