# modules/supernova/files/create.cron
# m h dom mon dow user   command
  0 8 *   *   *   create /opt/packages/create/create.py

