# modules/supernova/files/create.cron
# m h dom mon dow user   command
  15 16 *   *   *   create /opt/ocf/packages/create/create.py --keytab=/opt/create/private/create.keytab

