# modules/ocf/files/local/riot/crontab
# m h dom mon dow user command
*/15 * * * *      root /usr/bin/wget http://www.ocf.berkeley.edu/~staff/motd -O /opt/kiosk/motd -o /dev/null
