# modules/ocf_irc/files/inspircd.service.d/override.conf
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
Group=irc
Restart=always
