# modules/ocf_mesos/templates/master/marathon/marathon.service.erb
[Unit]
Description=Marathon
After=network.target
Wants=network.target

[Service]
EnvironmentFile=/opt/share/mesos/master/marathon-environ
# mesos_leader_ui_url to allow viewing stdout/stderr (rt#5645)
# reconciliation_interval to speed up recovery during restarts (rt#5651)
ExecStart=/usr/bin/marathon \
    --master zk://<%= @zookeeper_host %>/mesos \
    --zk zk://<%= @zookeeper_host %>/marathon \
    --hostname <%= @marathon_hostname %> \
    --mesos_authentication \
    --mesos_authentication_principal ocf \
    --mesos_authentication_secret_file /opt/share/mesos/master/marathon-secret \
    --mesos_leader_ui_url https://mesos.ocf.berkeley.edu/ \
    --reconciliation_interval 30000
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
