# modules/ocf_mirrors/templates/systemd/monitor.service.erb
[Unit]
Description=Healthcheck service for <%= @title %>
StartLimitInterval=70
StartLimitBurst=5

[Service]
User=<%= @exec_user %>
Type=simple
<% @environments.each do |env, val| -%>
Environment="<%= env %>=<%= val %>"
<% end -%>
ExecStart=<%= @exec_start %>
RuntimeMaxSec=10
