# modules/ocf_mirrors/templates/systemd/sync.timer.erb
[Unit]
Description=<% if @type == "sync" %>Sync<% elsif @type == "monitor" %>Monitoring<% end %> timer for <%= @title %>
After=network.target
Wants=network.target

[Timer]
OnCalendar=<%= @on_calendar %>

[Install]
WantedBy=timers.target
