# modules/ocf_mirrors/templates/systemd/sync.timer.erb
[Unit]
Description=Sync timer for <%= @title %>
After=network.target
Wants=network.target

[Timer]
OnCalendar=<%= @on_calendar %>

[Install]
WantedBy=timers.target
