# modules/ocf_mirrors/templates/systemd/sync.service.erb
[Unit]
Description=Sync service for <%= @title %>

[Service]
User=<%= @exec_user %>
Type=simple
<% @environments.each do |env, val| -%>
Environment="<%= env %>=<%= val %>"
<% end -%>
ExecStart=<%= @exec_start %>
# ~~one day should be good enough for anyone~~
# not enough for opensuse, it seems
RuntimeMaxSec=172800
